# laravel/workers.ini.default
;; laravel-workers.ini
;; https://laravel.com/docs/5.5/queues#supervisor-configuration
[program:laravel-worker]
process_name=worker_%(process_num)02d
directory=/var/www/html
command=php artisan queue:work redis --sleep=3
autostart=true
autorestart=true
user=www-data
numprocs=5
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
