# supervisor.d/crond.ini
[program:crond]
command=/usr/sbin/crond -f -d 8
autostart=true
autorestart=true
priority=30
user=root
numprocs=1
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
