# supervisor.d/crond.ini
[program:crond]
command=/usr/sbin/crond -f -S
autostart=true
autorestart=true
priority=30
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisord/crond.log
