# supervisor.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g "daemon off; pid /var/run/nginx.pid;"
autostart=true
autorestart=true
priority=31
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisord/nginx.log
