# supervisor.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm -F
autostart=true
autorestart=true
priority=50
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisord/php-fpm.log
