# supervisor.d/syslogd.ini
[program:syslogd]
command=/sbin/syslogd -n -O -
autostart=true
autorestart=true
priority=10
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisord/syslogd.log
