# supervisor.d/syslogd.ini
[program:syslogd]
command=/sbin/syslogd -n -O -
autostart=true
autorestart=true
priority=10
user=root
numprocs=1
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
